# meta-mender-demo/recipes-core/systemd/files/blinkenlights.service
[Unit]
Description=Blink LED2

[Service]
ExecStart=/usr/sbin/blinkenlights.sh

[Install]
WantedBy=multi-user.target
